# settings-portal/io.elementary.settings-daemon.xdg-desktop-portal.service.in
[Unit]
Description=Portal service (Pantheon Settings Daemon)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.elementary.settings-daemon
ExecStart=@libexecdir@/io.elementary.settings-daemon.xdg-desktop-portal
