# data/io.elementary.settings-daemon.system-update.timer
[Unit]
Description=Check for system updates daily

[Timer]
OnCalendar=daily
OnStartupSec=60
Persistent=true

[Install]
WantedBy=timers.target
